# quiet-edge.service
[Unit]
Description=quiet-edge iDRAC Fan Controller
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/quiet_edge_fan_control.py
Restart=always
RestartSec=5
# Failsafe: restore iDRAC dynamic auto control on stop/crash/shutdown
ExecStopPost=/usr/bin/ipmitool raw 0x30 0x30 0x01 0x01

[Install]
WantedBy=multi-user.target